# .config/autostart/plasma_startup.desktop
[Desktop Entry]
Exec=/home/nacl10/scripts/plasma_startup
Icon=dialog-scripts
Name=plasma_startup
Path=
Type=Application
X-KDE-AutostartScript=true
